# debian/synapse/supervisord.conf
[supervisord]
nodaemon=true

[program:x11vnc]
command=x11vnc -forever -rfbport 5900
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

